# biasimpacter/crontab
* * * * * python /biasimpacter/biasimpacter/app/app.py >> /biasimpacter/logs/crontab.log 2>&1
